# crontab
# Run every day at 09:00 UTC
# Min Hour DayOfMonth Month DayOfWeek
0 7 * * * . /app/cron.env && /usr/local/bin/python /app/post_image.py >> /var/log/cron.log 2>&1
